# crontab
# 정해진 시간에 epg2xml을 실행하여 xmltv.xml 파일로 출력
56 4,16 * * * cd /frontend/epg && /usr/bin/python3 -m epg2xml run --xmlfile=/frontend/epg/xmltv.xml >> /var/log/cron.log 2>&1
