# crontab
# 매일 새벽 2시에 epg2xml을 실행하여 xmltv.xml 파일로 출력
0 2 * * * /usr/bin/python3 -m epg2xml run --xmlfile=/app/xmltv.xml >> /var/log/cron.log 2>&1
